# apps/tcdevices/settings/supervisor.d/tcdevices.ini
[program:tcdevices]
command=/home/super/.local/share/virtualenvs/tcloudservers-IY81mg_K/bin/python run.py
stdout_logfile=/home/super/logs/tcloud/supervisord-cmdb-out.log
stderr_logfile=/home/super/logs/tcloud/supervisord-cmdb-error.log
directory=/home/super/py3/tcloudservers/apps/tcdevices
autostart=true
autorestart=true
startsecs=10
startretries=3
